# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ==========================================================================
; Core Platform
; ==========================================================================

[program:health]
command=/app/bin/health
environment=PORT=%(ENV_HEALTH_PORT)s
directory=/app
stdout_logfile=/app/logs/health.log
stderr_logfile=/app/logs/health.err
autorestart=true
startretries=5
priority=10
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:notification]
command=/app/bin/notification
environment=PORT=%(ENV_NOTIFICATION_PORT)s
directory=/app
stdout_logfile=/app/logs/notification.log
stderr_logfile=/app/logs/notification.err
autorestart=true
startretries=5
priority=10
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:auth]
command=/app/bin/auth
environment=PORT=%(ENV_AUTH_PORT)s
directory=/app
stdout_logfile=/app/logs/auth.log
stderr_logfile=/app/logs/auth.err
autorestart=true
startretries=5
priority=20
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:entitlement]
command=/app/bin/entitlement
environment=PORT=%(ENV_ENTITLEMENT_PORT)s
directory=/app
stdout_logfile=/app/logs/entitlement.log
stderr_logfile=/app/logs/entitlement.err
autorestart=true
startretries=5
priority=20
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:billing]
command=/app/bin/billing
environment=PORT=%(ENV_BILLING_PORT)s
directory=/app
stdout_logfile=/app/logs/billing.log
stderr_logfile=/app/logs/billing.err
autorestart=true
startretries=5
priority=20
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:backend]
command=/app/bin/backend
environment=PORT=%(ENV_BACKEND_PORT)s
directory=/app
stdout_logfile=/app/logs/backend.log
stderr_logfile=/app/logs/backend.err
autorestart=true
startretries=5
priority=30
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:frontend]
command=npx next start --port %(ENV_FRONTEND_PORT)s
directory=/app/frontend
environment=NODE_ENV=production
stdout_logfile=/app/logs/frontend.log
stderr_logfile=/app/logs/frontend.err
autorestart=true
priority=30
stopsignal=SIGTERM
stopwaitsecs=20
stopasgroup=true
killasgroup=true

[program:storage]
command=/app/bin/storage
environment=PORT=%(ENV_STORAGE_PORT)s
directory=/app
stdout_logfile=/app/logs/storage.log
stderr_logfile=/app/logs/storage.err
autorestart=true
priority=30
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:cloud-adapter]
command=/app/bin/cloud-adapter
environment=PORT=%(ENV_CLOUD_ADAPTER_PORT)s
directory=/app
stdout_logfile=/app/logs/cloud-adapter.log
stderr_logfile=/app/logs/cloud-adapter.err
autorestart=true
priority=30
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

; ==========================================================================
; Operations & Job Workers
; ==========================================================================

[program:scheduler]
command=/app/bin/scheduler
environment=PORT=0
directory=/app
stdout_logfile=/app/logs/scheduler.log
stderr_logfile=/app/logs/scheduler.err
autorestart=true
priority=40
stopsignal=SIGTERM
stopwaitsecs=15
stopasgroup=true
killasgroup=true

[program:inventory]
command=/app/bin/inventory
environment=PORT=%(ENV_INVENTORY_PORT)s
directory=/app
stdout_logfile=/app/logs/inventory.log
stderr_logfile=/app/logs/inventory.err
autorestart=true
priority=40
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:orders]
command=/app/bin/orders
environment=PORT=%(ENV_ORDERS_PORT)s
directory=/app
stdout_logfile=/app/logs/orders.log
stderr_logfile=/app/logs/orders.err
autorestart=true
priority=40
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

; ==========================================================================
; Communication & Payment Plugins
; ==========================================================================

[program:email]
command=/app/bin/email
environment=PORT=%(ENV_EMAIL_PORT)s
directory=/app
stdout_logfile=/app/logs/email.log
stderr_logfile=/app/logs/email.err
autorestart=true
priority=50
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:sms]
command=/app/bin/sms
environment=PORT=%(ENV_SMS_PORT)s
directory=/app
stdout_logfile=/app/logs/sms.log
stderr_logfile=/app/logs/sms.err
autorestart=true
priority=50
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:whatsapp]
command=/app/bin/whatsapp
environment=PORT=%(ENV_WHATSAPP_PORT)s
directory=/app
stdout_logfile=/app/logs/whatsapp.log
stderr_logfile=/app/logs/whatsapp.err
autorestart=true
priority=50
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:telegram]
command=/app/bin/telegram
environment=PORT=%(ENV_TELEGRAM_PORT)s
directory=/app
stdout_logfile=/app/logs/telegram.log
stderr_logfile=/app/logs/telegram.err
autorestart=true
priority=50
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:razorpay]
command=/app/bin/razorpay
environment=PORT=%(ENV_RAZORPAY_PORT)s
directory=/app/env/razorpay
stdout_logfile=/app/logs/razorpay.log
stderr_logfile=/app/logs/razorpay.err
autorestart=true
priority=50
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:stripe]
command=/app/bin/stripe
environment=PORT=%(ENV_STRIPE_PORT)s
directory=/app/env/stripe
stdout_logfile=/app/logs/stripe.log
stderr_logfile=/app/logs/stripe.err
autorestart=true
priority=50
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

; ==========================================================================
; Astrology Plugins (Node.js)
; ==========================================================================

[program:jhora]
command=node src/index.js
environment=PORT=%(ENV_JHORA_PORT)s
directory=/app/jhora
stdout_logfile=/app/logs/jhora.log
stderr_logfile=/app/logs/jhora.err
autorestart=true
priority=50
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:vedika]
command=node src/index.js
environment=PORT=%(ENV_VEDIKA_PORT)s
directory=/app/vedika
stdout_logfile=/app/logs/vedika.log
stderr_logfile=/app/logs/vedika.err
autorestart=true
priority=50
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
