# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ==========================================================================
; Core Platform
; ==========================================================================

[program:backend]
command=/app/bin/backend
environment=PORT=%(ENV_BACKEND_PORT)s
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=10

[program:frontend]
command=npx next start --port %(ENV_FRONTEND_PORT)s
directory=/app/frontend
environment=NODE_ENV=production
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=20

[program:scheduler]
command=/app/bin/scheduler
environment=PORT=0
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=20

; ==========================================================================
; Microservices
; ==========================================================================

[program:entitlement]
command=/app/bin/entitlement
environment=PORT=%(ENV_ENTITLEMENT_PORT)s
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=30

[program:billing]
command=/app/bin/billing
environment=PORT=%(ENV_BILLING_PORT)s
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=30

[program:health]
command=/app/bin/health
environment=PORT=%(ENV_HEALTH_PORT)s
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=30

[program:marketing]
command=/app/bin/marketing
environment=PORT=%(ENV_MARKETING_PORT)s
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=30

[program:cloud-adapter]
command=/app/bin/cloud-adapter
environment=PORT=%(ENV_CLOUD_ADAPTER_PORT)s
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=30

[program:notification]
command=/app/bin/notification
environment=PORT=%(ENV_NOTIFICATION_PORT)s
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=30

[program:storage]
command=/app/bin/storage
environment=PORT=%(ENV_STORAGE_PORT)s
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=30

[program:inventory]
command=/app/bin/inventory
environment=PORT=%(ENV_INVENTORY_PORT)s
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=30

[program:orders]
command=/app/bin/orders
environment=PORT=%(ENV_ORDERS_PORT)s
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=30

; ==========================================================================
; Communication Plugins
; ==========================================================================

[program:email]
command=/app/bin/email
environment=PORT=%(ENV_EMAIL_PORT)s
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=40

[program:sms]
command=/app/bin/sms
environment=PORT=%(ENV_SMS_PORT)s
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=40

[program:whatsapp]
command=/app/bin/whatsapp
environment=PORT=%(ENV_WHATSAPP_PORT)s
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=40

[program:telegram]
command=/app/bin/telegram
environment=PORT=%(ENV_TELEGRAM_PORT)s
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=40

; ==========================================================================
; Payment Plugins
; ==========================================================================

[program:razorpay]
command=/app/bin/razorpay
environment=PORT=%(ENV_RAZORPAY_PORT)s
directory=/app/env/razorpay
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=40

[program:stripe]
command=/app/bin/stripe
environment=PORT=%(ENV_STRIPE_PORT)s
directory=/app/env/stripe
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=40

; ==========================================================================
; Astrology Plugins (Node.js)
; ==========================================================================

[program:jhora]
command=node src/index.js
environment=PORT=%(ENV_JHORA_PORT)s
directory=/app/jhora
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=40

[program:vedika]
command=node src/index.js
environment=PORT=%(ENV_VEDIKA_PORT)s
directory=/app/vedika
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=40
