# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ==========================================================================
; Core Platform
; ==========================================================================

[program:backend]
command=/app/bin/backend
environment=PORT=8080
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=10

[program:frontend]
command=npx next start --port 3000
directory=/app/frontend
environment=NODE_ENV=production
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=20

[program:scheduler]
command=/app/bin/scheduler
environment=PORT=0
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=20

; ==========================================================================
; Microservices
; ==========================================================================

[program:entitlement]
command=/app/bin/entitlement
environment=PORT=8085
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=30

[program:billing]
command=/app/bin/billing
environment=PORT=8090
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=30

[program:health]
command=/app/bin/health
environment=PORT=8091
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=30

[program:marketing]
command=/app/bin/marketing
environment=PORT=8092
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=30

[program:cloud-adapter]
command=/app/bin/cloud-adapter
environment=PORT=8093
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=30

[program:notification]
command=/app/bin/notification
environment=PORT=8094
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=30

[program:storage]
command=/app/bin/storage
environment=PORT=8098
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=30

[program:inventory]
command=/app/bin/inventory
environment=PORT=8096
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=30

[program:orders]
command=/app/bin/orders
environment=PORT=8097
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=30

; ==========================================================================
; Communication Plugins
; ==========================================================================

[program:email]
command=/app/bin/email
environment=PORT=8081
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=40

[program:sms]
command=/app/bin/sms
environment=PORT=8082
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=40

[program:whatsapp]
command=/app/bin/whatsapp
environment=PORT=8083
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=40

[program:telegram]
command=/app/bin/telegram
environment=PORT=8084
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=40

; ==========================================================================
; Payment Plugins
; ==========================================================================

[program:razorpay]
command=/app/bin/razorpay
environment=PORT=8086
directory=/app/env/razorpay
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=40

[program:stripe]
command=/app/bin/stripe
environment=PORT=8087
directory=/app/env/stripe
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=40

; ==========================================================================
; Astrology Plugins (Node.js)
; ==========================================================================

[program:jhora]
command=node src/index.js
environment=PORT=3130
directory=/app/jhora
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=40

[program:vedika]
command=node src/index.js
environment=PORT=3140
directory=/app/vedika
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=40
