# docker/supervisord-debug.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ==========================================================================
; Core Platform (Go services with Delve)
; ==========================================================================

[program:health]
command=dlv exec --headless --listen=0.0.0.0:%(ENV_DEBUG_HEALTH_PORT)s --api-version=2 -- /app/bin/health
environment=PORT=%(ENV_HEALTH_PORT)s
directory=/app
stdout_logfile=/app/logs/health.log
stderr_logfile=/app/logs/health.err
autorestart=true
startretries=5
priority=10
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:notification]
command=dlv exec --headless --listen=0.0.0.0:%(ENV_DEBUG_NOTIFICATION_PORT)s --api-version=2 -- /app/bin/notification
environment=PORT=%(ENV_NOTIFICATION_PORT)s
directory=/app
stdout_logfile=/app/logs/notification.log
stderr_logfile=/app/logs/notification.err
autorestart=true
startretries=5
priority=10
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:auth]
command=dlv exec --headless --listen=0.0.0.0:%(ENV_DEBUG_AUTH_PORT)s --api-version=2 -- /app/bin/auth
environment=PORT=%(ENV_AUTH_PORT)s
directory=/app
stdout_logfile=/app/logs/auth.log
stderr_logfile=/app/logs/auth.err
autorestart=true
startretries=5
priority=20
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:entitlement]
command=dlv exec --headless --listen=0.0.0.0:%(ENV_DEBUG_ENTITLEMENT_PORT)s --api-version=2 -- /app/bin/entitlement
environment=PORT=%(ENV_ENTITLEMENT_PORT)s
directory=/app
stdout_logfile=/app/logs/entitlement.log
stderr_logfile=/app/logs/entitlement.err
autorestart=true
startretries=5
priority=20
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:billing]
command=dlv exec --headless --listen=0.0.0.0:%(ENV_DEBUG_BILLING_PORT)s --api-version=2 -- /app/bin/billing
environment=PORT=%(ENV_BILLING_PORT)s
directory=/app
stdout_logfile=/app/logs/billing.log
stderr_logfile=/app/logs/billing.err
autorestart=true
startretries=5
priority=20
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:backend]
command=dlv exec --headless --listen=0.0.0.0:%(ENV_DEBUG_BACKEND_PORT)s --api-version=2 -- /app/bin/backend
environment=PORT=%(ENV_BACKEND_PORT)s
directory=/app
stdout_logfile=/app/logs/backend.log
stderr_logfile=/app/logs/backend.err
autorestart=true
startretries=5
priority=30
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:frontend]
command=node --inspect=0.0.0.0:9230 ./node_modules/.bin/next dev --port %(ENV_FRONTEND_PORT)s
directory=/app/frontend
environment=NODE_ENV=development
stdout_logfile=/app/logs/frontend.log
stderr_logfile=/app/logs/frontend.err
autorestart=true
priority=30
stopsignal=SIGTERM
stopwaitsecs=20
stopasgroup=true
killasgroup=true

[program:storage]
command=dlv exec --headless --listen=0.0.0.0:%(ENV_DEBUG_STORAGE_PORT)s --api-version=2 -- /app/bin/storage
environment=PORT=%(ENV_STORAGE_PORT)s
directory=/app
stdout_logfile=/app/logs/storage.log
stderr_logfile=/app/logs/storage.err
autorestart=true
priority=30
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:cloud-adapter]
command=dlv exec --headless --listen=0.0.0.0:%(ENV_DEBUG_CLOUD_ADAPTER_PORT)s --api-version=2 -- /app/bin/cloud-adapter
environment=PORT=%(ENV_CLOUD_ADAPTER_PORT)s
directory=/app
stdout_logfile=/app/logs/cloud-adapter.log
stderr_logfile=/app/logs/cloud-adapter.err
autorestart=true
priority=30
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

; ==========================================================================
; Operations & Job Workers
; ==========================================================================

[program:scheduler]
command=dlv exec --headless --listen=0.0.0.0:%(ENV_DEBUG_SCHEDULER_PORT)s --api-version=2 -- /app/bin/scheduler
environment=PORT=0
directory=/app
stdout_logfile=/app/logs/scheduler.log
stderr_logfile=/app/logs/scheduler.err
autorestart=true
priority=40
stopsignal=SIGTERM
stopwaitsecs=15
stopasgroup=true
killasgroup=true

[program:inventory]
command=dlv exec --headless --listen=0.0.0.0:%(ENV_DEBUG_INVENTORY_PORT)s --api-version=2 -- /app/bin/inventory
environment=PORT=%(ENV_INVENTORY_PORT)s
directory=/app
stdout_logfile=/app/logs/inventory.log
stderr_logfile=/app/logs/inventory.err
autorestart=true
priority=40
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:orders]
command=dlv exec --headless --listen=0.0.0.0:%(ENV_DEBUG_ORDERS_PORT)s --api-version=2 -- /app/bin/orders
environment=PORT=%(ENV_ORDERS_PORT)s
directory=/app
stdout_logfile=/app/logs/orders.log
stderr_logfile=/app/logs/orders.err
autorestart=true
priority=40
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

; ==========================================================================
; Communication & Payment Plugins
; ==========================================================================

[program:email]
command=dlv exec --headless --listen=0.0.0.0:%(ENV_DEBUG_EMAIL_PORT)s --api-version=2 -- /app/bin/email
environment=PORT=%(ENV_EMAIL_PORT)s
directory=/app
stdout_logfile=/app/logs/email.log
stderr_logfile=/app/logs/email.err
autorestart=true
priority=50
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:sms]
command=dlv exec --headless --listen=0.0.0.0:%(ENV_DEBUG_SMS_PORT)s --api-version=2 -- /app/bin/sms
environment=PORT=%(ENV_SMS_PORT)s
directory=/app
stdout_logfile=/app/logs/sms.log
stderr_logfile=/app/logs/sms.err
autorestart=true
priority=50
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:whatsapp]
command=dlv exec --headless --listen=0.0.0.0:%(ENV_DEBUG_WHATSAPP_PORT)s --api-version=2 -- /app/bin/whatsapp
environment=PORT=%(ENV_WHATSAPP_PORT)s
directory=/app
stdout_logfile=/app/logs/whatsapp.log
stderr_logfile=/app/logs/whatsapp.err
autorestart=true
priority=50
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:telegram]
command=dlv exec --headless --listen=0.0.0.0:%(ENV_DEBUG_TELEGRAM_PORT)s --api-version=2 -- /app/bin/telegram
environment=PORT=%(ENV_TELEGRAM_PORT)s
directory=/app
stdout_logfile=/app/logs/telegram.log
stderr_logfile=/app/logs/telegram.err
autorestart=true
priority=50
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:razorpay]
command=dlv exec --headless --listen=0.0.0.0:%(ENV_DEBUG_RAZORPAY_PORT)s --api-version=2 -- /app/bin/razorpay
environment=PORT=%(ENV_RAZORPAY_PORT)s
directory=/app/env/razorpay
stdout_logfile=/app/logs/razorpay.log
stderr_logfile=/app/logs/razorpay.err
autorestart=true
priority=50
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:stripe]
command=dlv exec --headless --listen=0.0.0.0:%(ENV_DEBUG_STRIPE_PORT)s --api-version=2 -- /app/bin/stripe
environment=PORT=%(ENV_STRIPE_PORT)s
directory=/app/env/stripe
stdout_logfile=/app/logs/stripe.log
stderr_logfile=/app/logs/stripe.err
autorestart=true
priority=50
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

; ==========================================================================
; Astrology Plugins (Node.js with --inspect)
; ==========================================================================

[program:jhora]
command=node --inspect=0.0.0.0:9231 src/index.js
environment=PORT=%(ENV_JHORA_PORT)s
directory=/app/jhora
stdout_logfile=/app/logs/jhora.log
stderr_logfile=/app/logs/jhora.err
autorestart=true
priority=50
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true

[program:vedika]
command=node --inspect=0.0.0.0:9232 src/index.js
environment=PORT=%(ENV_VEDIKA_PORT)s
directory=/app/vedika
stdout_logfile=/app/logs/vedika.log
stderr_logfile=/app/logs/vedika.err
autorestart=true
priority=50
stopsignal=SIGTERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
